# site/minecraft/templates/minecraft.service.epp
[Unit]
Description = Minecraft Server

Wants=network.target
After=netwrok.target

[Service]
WorkingDirectory=<%= $install_dir %>
ExecStart=/usr/bin/java -Xmx512M -Xms32M -jar server.jar nogui

[Install]
WantedBy=multiuser.target
